# launchd/com.fleet.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Fleet Supervisor LaunchAgent
  NOTE: This is a TEMPLATE. install.sh will generate the real plist
  with $HOME paths resolved for your user. Do not load this file directly.
-->
<dict>
    <key>Label</key>
    <string>com.fleet.supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HANDLER_DIR__/fleet-supervisor.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude-fleet/logs/supervisor-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude-fleet/logs/supervisor-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
